# roles/setdefaultpath/template/etc/profile.d/path.sh
export PATH=/opt/node-v16.19.1-linux-x64/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin

export JAVA_HOME=/usr/lib/jvm/java-11-openjdk-amd64
export CLASSPATH="/usr/lib/jvm/java-11-openjdk-amd64/lib:."
