# roles/wol_servicio/template/wol.service.j2
[Unit]
Description=Configure Wake On LAN
After=network.targe

[Service]
Type=oneshot
ExecStart=/sbin/ethtool -s {{ ansible_interfaces[0] }} wol g
ReaminAfterExit=yes

[Install]
WantedBy=basic.target
